# model_server/mac_services/com.projectxylen.datacollector.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.projectxylen.datacollector</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/opt/trading_model/venv/bin/python</string>
        <string>/usr/local/opt/trading_model/data_collector.py</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/usr/local/opt/trading_model</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/opt/trading_model/venv/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/usr/local/opt/trading_model/logs/data_collector.log</string>
    
    <key>StandardErrorPath</key>
    <string>/usr/local/opt/trading_model/logs/data_collector.error.log</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>ProcessType</key>
    <string>Background</string>
    
    <key>Nice</key>
    <integer>5</integer>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
